# pr.service
[Unit]
Description=PR

[Service]
ExecStart=/usr/local/bin/pr

[Install]
WantedBy=multi-user.target
