# pr.service
[Unit]
Description=PR

[Service]
ExecStart=/usr/sbin/pr

[Install]
WantedBy=multi-user.target
